# src/main/resources/org/jenkinsci/modules/upstart_slave_installer/jenkins-slave.conf
# send output to /var/log/upstart/
console log

description "Jenkins slave"
exec java -jar {jar} {args}
respawn
setuid {username}

start on runlevel [2345]
stop on runlevel [016]

